# AUR/novpn.install
post_install() {
    novpn_ns configure_default
    systemctl daemon-reload
    systemctl enable novpn.service
    systemctl start novpn.service
}

pre_upgrade() {
    novpn_ns down
    systemctl stop novpn.service
    systemctl disable novpn.service
}

post_upgrade() {
    systemctl daemon-reload
    systemctl enable novpn.service
    systemctl start novpn.service
}


pre_remove() {
    novpn_ns down
    systemctl stop novpn.service
    systemctl disable novpn.service
}

post_remove() {
    systemctl daemon-reload
}
